# DOCKER_SERVICES/SUBSET/RADIUS_MYSQL_XRDP_SSH/RADIUS_FINAL/supervisord.conf
[supervisord]
nodaemon=true

[program:freeradius]
command= freeradius -X

[eventlistener:processes]
events= PROCESS_STATE_EXITED
command= /etc/supervisor/conf.d/stop-supervisor.sh


