# DOCKER_SERVICES/MYIMAGES/VULNNETWORK/XRDP_FINAL/supervisord.conf
[supervisord]
nodaemon=true

[program:rdp]
command= xrdp -nodaemon

[eventlistener:processes]
events= PROCESS_STATE_EXITED
command= /etc/supervisor/conf.d/stop-supervisor.sh


