# DOCKER_SERVICES/MYIMAGES/VULNWEB/TRAC_FINAL/supervisord.conf
[supervisord]
nodaemon=true

[program:trac]
command= tracd --port 8000 /var/lib/trac/test

[eventlistener:processes]
events= PROCESS_STATE_EXITED
command= /etc/supervisor/conf.d/stop-supervisor.sh


