# HONEYPOT/FTPPROXYARCH/supervisord.conf
[supervisord]
nodaemon=true





[program:ftp]
command=python /PROXIES/FTP_PROXY.py

[eventlistener:processes]
events=PROCESS_STATE_EXITED
command=/PROXIES/stop-supervisor.sh



