# HONEYPOT/SSHPROXYARCH/supervisord.conf
[supervisord]
nodaemon=true





[program:ssh]
command=python /PROXIES/Mod_sshpot.py

[eventlistener:processes]
events=PROCESS_STATE_EXITED
command=/PROXIES/stop-supervisor.sh


