# DOCKER_SERVICES/SUBSET/HTTP_MYSQL_SMTP_SSH/http-subset/VULNFTP/FTP_FINAL/supervisord.conf
[supervisord]
nodaemon=true

[program:proftpd]
command= /usr/local/bin/start.sh

[eventlistener:processes]
events= PROCESS_STATE_EXITED
command= /etc/supervisor/conf.d/stop-supervisor.sh


